# build/mysql56/files/manifest-mysql-56.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type="manifest" name="MySQL">
    <service
            name="application/database/mysql_56"
            type="service"
            version="1">

        <!-- Wait for all local filesystems to be mounted -->
        <dependency
                name="filesystem"
                grouping="require_all"
                restart_on="none"
                type="service">
            <service_fmri value="svc:/system/filesystem/local" />
        </dependency>

        <!-- Wati for network interfaces to be initialized -->
        <dependency
                name='network'
                grouping='require_all'
                restart_on='none'
                type='service'>
                <service_fmri value='svc:/milestone/network:default' />
        </dependency>


        <!-- Execution method for start, stop and restart -->
        <exec_method
                type="method"
                name="start"
                exec="/lib/svc/method/mysql_56 start"
                timeout_seconds="30" />

        <exec_method
                type="method"
                name="stop"
                exec="/lib/svc/method/mysql_56 stop"
                timeout_seconds="30" />

        <exec_method
                type="method"
                name="refresh"
                exec="/lib/svc/method/mysql_56 reload"
                timeout_seconds="30" />

        <!-- Other properties and configs -->
        <property_group name="general" type="framework">
            <propval name="value_authorization" type="astring"
                    value="solaris.smf.value.mysql" />
            <propval name="action_authorization" type="astring"
                    value="solaris.smf.manage.mysql" />
        </property_group>

        <instance name="default" enabled="false">
            <method_context>
                <method_credential user='mysql' group='mysql'/>
            </method_context>

            <property_group name='mysql_56' type='application'>
                <propval name='base' type='astring'
                    value='/usr/local/mysql/5.6'/>
                <propval name='bin' type='astring'
                    value='/usr/local/mysql/5.6/bin'/>
                <propval name='data' type='astring'
                    value='/var/mysql/5.6/data'/>
                <propval name='mycnf' type='astring'
                    value='/var/mysql/5.6/etc/my.cnf'/>
                <propval name='socket' type='astring'
                    value='/var/mysql/5.6/run/mysqld.socket'/>
                <propval name='pidfile' type='astring'
                    value='/var/mysql/5.6/run/mysqld.pid'/>
            </property_group>
        </instance>

        <stability value="Evolving" />

        <template>
            <common_name>
                <loctext xml:lang="C">MySQL Database Server version 5.6</loctext>
            </common_name>

            <!-- Define documentation path (MAN-page path) -->
            <documentation>
                <manpage title="mysql" section="1" manpath="/usr/local/mysql/5.6/man" />
                <manpage title="mysqld" section="1" manpath="/usr/local/mysql/5.6/man" />
                <manpage title="mysqladmin" section="1" manpath="/usr/local/mysql/5.6/man" />
                <doc_link name="dev.mysql.org" uri="http://dev.mysql.org" />
            </documentation>
        </template>

    </service>
</service_bundle>

