# build/samba4/files/manifest-samba-smbd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="samba-smbd">

    <service name="network/samba/smbd" type="service" version="1">

        <create_default_instance enabled="false" />

        <single_instance/>

        <dependency name="net-loopback" grouping="require_any" restart_on="none" type="service">
            <service_fmri value="svc:/network/loopback" />
        </dependency>

        <dependency name="net-service" grouping="require_all" restart_on="none" type="service">
            <service_fmri value="svc:/network/service"/>
        </dependency>

        <dependency name="net-physical" grouping="require_all" restart_on="none" type="service">
            <service_fmri value="svc:/network/physical" />
        </dependency>

        <dependency name="filesystem-local" grouping="require_all" restart_on="none" type="service">
            <service_fmri value="svc:/system/filesystem/local" />
        </dependency>

        <dependent name="samba-smbd-multi-user-server" grouping="optional_all" restart_on="none">
            <service_fmri value="svc:/milestone/multi-user-server" />
        </dependent>

        <exec_method type="method" name="start" exec="/usr/local/sbin/smbd -D" timeout_seconds="170" />

        <exec_method type="method" name="stop" exec="/usr/bin/kill `cat /var/samba/locks/smbd.pid`" timeout_seconds="60" />

        <stability value="Unstable" />

        <template>
            <common_name>
                <loctext xml:lang="C">
                SMB file server - smbd daemon
                </loctext>
            </common_name>
            <documentation>
                <manpage title="smbd" section="1m" manpath="/usr/local/man" />
                <manpage title="smb.conf" section="4" manpath="/usr/local/man" />
            </documentation>
        </template>
    </service>

</service_bundle>
