# build/neo4j/files/manifest-neo4j.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
  Copyright 2004 Sun Microsystems, Inc.  All rights reserved.
  Use is subject to license terms.
-->

<service_bundle type='manifest' name='neo4j'>

<service name='database/neo4j'
         type='service'
         version='1'>

  <single_instance />

  <dependency name='fs-local'
              grouping='require_all'
              restart_on='none'
              type='service'>
    <service_fmri value='svc:/system/filesystem/local' />
  </dependency>

  <dependency name='network-service'
              grouping='require_all'
              restart_on='none'
              type='service'>
    <service_fmri value='svc:/network/service' />
  </dependency>

  <dependency name='name-services'
              grouping='require_all'
              restart_on='refresh'
              type='service'>
    <service_fmri value='svc:/milestone/name-services' />
  </dependency>

  <dependency name='identity'
              grouping='optional_all'
              restart_on='refresh'
              type='service'>
    <service_fmri value='svc:/system/identity:domain' />
  </dependency>

  <dependency name='system-log'
              grouping='optional_all'
              restart_on='none'
              type='service'>
    <service_fmri value='svc:/system/system-log' />
  </dependency>

  <exec_method type='method'
               name='start'
               exec='/usr/local/neo4j/bin/neo4j start'
               timeout_seconds='130'/>

  <exec_method type='method'
               name='stop'
               exec='/usr/local/neo4j/bin/neo4j stop'
               timeout_seconds='130' />

  <exec_method type='method'
               name='restart'
               exec='/usr/local/neo4j/bin/neo4j restart'
               timeout_seconds='130' />


  <instance name='default' enabled='false'>
    <method_context>
        <method_credential user='neo4j' group='neo4j' />
        <method_environment>
            <envvar name='NEO4J_HOME' value='/opt/neo4j'/>
            <envvar name='NEO4J_INSTANCE' value='/var/neo4j'/>
            <envvar name='JAVA_HOME' value='/opt/java/jdk8'/>
        </method_environment>
    </method_context>
  </instance>

  <stability value='Evolving' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Neo4j graph database</loctext>
    </common_name>
  </template>
</service>
</service_bundle>
